# package/nier.service
[Unit]
Description=nier
After=network.target

[Service]
LimitAS=infinity
LimitRSS=infinity
LimitCORE=infinity
LimitNOFILE=65535
Type=simple
ExecStart=/usr/bin/nier
ExecStop=/usr/bin/kill $MAINPID
Restart=always
RestartSec=30s
User=root

[Install]
WantedBy=multi-user.target
